# debian/postrm
#!/bin/bash
set -e

# remove pypiserver start script
update-rc.d -f pypiserver remove
