# debian/python-pypiserver-daemon/DEBIAN/prerm
#!/bin/sh
set -e

# stop pypiserver server
invoke-rc.d pypiserver stop
